# servo-gecko-try/files/servo-gecko-try.conf
exec /home/servo/servo-gecko-try/_venv/bin/servo-gecko-try

setuid servo
setgid servo

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

env HOME=/home/servo

chdir /home/servo/servo-gecko-try

