# buildbot/master/files/buildbot-github-listener.conf
exec /usr/local/bin/github_buildbot.py -p 9010 -m localhost:9001 --auth=change:{{ pillar['buildbot']['credentials']['change-pass'] }} --secret={{ pillar['buildbot']['credentials']['gh-hook-secret'] }} -l {{ common.servo_home }}/buildbot/master/github-listener.log

setuid servo
setgid servo

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

env HOME={{ common.servo_home }}

