# xvfb/xvfb.conf
exec Xvfb :0 -screen 0 1280x960x24 +extension RANDR +extension RENDER +extension GLX

start on runlevel [2345]
stop on runlevel [016]

