# buildbot/buildbot-slave.conf
exec /usr/local/bin/buildslave start --nodaemon /home/servo/buildbot/slave

setuid servo
setgid servo

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

env HOME=/home/servo

